# debian/wazo-websocketd.service
[Unit]
Description=wazo-websocketd server
ConditionPathExists=!/var/lib/wazo/disabled
After=network.target rabbitmq-server.service wazo-auth.service
StartLimitBurst=15
StartLimitIntervalSec=150

[Service]
ExecStart=/usr/bin/wazo-websocketd
ExecStartPost=/usr/bin/wazo-websocketd-wait-online
Restart=on-failure
RestartSec=5
SyslogIdentifier=wazo-websocketd

[Install]
WantedBy=multi-user.target
